# debian/assets/xbnet.service
[Unit]
Description=Example

[Service]
ExecStart=/usr/bin/xbnet

[Install]
WantedBy=multi-user.target